# Linux/etc/udev/rules.d/56-mouse-detach.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="046d", ATTR{idProduct}=="c05a", RUN+="/usr/local/bin/findxinput.sh"
